# deploy/void-wg.service
[Unit]
Description=void-wg WireGuard Control Plane
Documentation=https://github.com/wester11/void_wg
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/void-wg
EnvironmentFile=-/opt/void-wg/.env
ExecStartPre=-/usr/bin/docker compose pull --quiet
ExecStart=/usr/bin/docker compose up -d --build --remove-orphans
ExecStop=/usr/bin/docker compose down
TimeoutStartSec=600

[Install]
WantedBy=multi-user.target
